# crontab
# m h dom mon dow command
  * *  *   *    *  /home/joseph/Desktop/temp_logger/i2c_to_py_temp.sh

